# etc/supervisor/conf.d/hue.conf
[program:hue]
user=hue
command=/home/hue/hue-3.8.1/build/env/bin/hue runserver